# .config/systemd/user/bing-wallpaper.service
[Unit]
Description=Fetch Bing daily wallpaper

[Service]
Type=oneshot
ExecStart=/run/current-system/sw/bin/bash -c 'echo $PATH'
